# .config/systemd/user/mbsync.service
[Unit]
Description=Mailbox synchronization service

[Service]
Type=oneshot
ExecStart=/usr/bin/env mbsync -Va
ExecStartPost=/usr/bin/env notmuch new
